# software/redirect-server/supervisor-server.conf
[program:server]
command=python3 /redirect-server.py
